# deploy/scheduler.service
[Unit]
Description=In-app scheduler (threading-based) for arm_songs
After=network.target gunicorn.service

[Service]
User=arm
Group=www-data
WorkingDirectory=/srv/arm_songs
EnvironmentFile=/srv/arm_songs/.env
ExecStart=/srv/arm_songs/.venv/bin/python manage.py run_scheduler --minutes 30
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
